# services/deploy/systemd/diabetes-bot.service.example
## Example unit file for running the Diabetes bot
## Copy to /etc/systemd/system/diabetes-bot.service and adjust paths and user

[Unit]
Description=Diabetes Telegram bot
After=network.target

# Service configuration
[Service]
Type=simple
User=www-data
PermissionsStartOnly=true
WorkingDirectory=/srv/diabetes-bot
EnvironmentFile=/srv/diabetes-bot/.env
Environment=MPLCONFIGDIR=/srv/diabetes-bot/data/mpl-cache
ExecStartPre=/usr/bin/install -d -o www-data -g www-data -m 700 /srv/diabetes-bot/data/mpl-cache
ExecStartPre=/usr/bin/curl -fsS "$PUBLIC_ORIGIN/health"
ExecStart=/srv/diabetes-bot/scripts/run_bot.sh
Restart=on-failure
RestartSec=5
# log to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
