# docs/deploy/diabetes-bot.service
[Unit]
Description=Diabetes Telegram bot
After=network.target

[Service]
Type=simple
User=www-data
StateDirectory=diabetes-bot
WorkingDirectory=/opt/saharlight-ux
EnvironmentFile=/opt/saharlight-ux/.env
Environment=MPLCONFIGDIR=/opt/saharlight-ux/data/mpl-cache
ExecStartPre=/usr/bin/install -d -o www-data -g www-data -m 700 /opt/saharlight-ux/data/mpl-cache
ExecStart=/opt/saharlight-ux/scripts/run_bot.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
