# docs/deploy/diabetes-assistant.service
[Unit]
Description=Diabetes Assistant WebApp
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/diabetes-assistant
EnvironmentFile=/opt/diabetes-assistant/.env
ExecStart=/usr/bin/env uvicorn webapp.server:app --host 0.0.0.0 --port 8000 --workers 4
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
