# docs/deploy/diabetes-bot.service
[Unit]
Description=Diabetes Telegram bot
After=network.target

[Service]
Type=simple
User=www-data
# systemd will create /var/lib/diabetes-bot and expose it via STATE_DIRECTORY
StateDirectory=diabetes-bot
# use systemd-managed cache directory for matplotlib config
CacheDirectory=diabetes-bot
WorkingDirectory=/opt/saharlight-ux
EnvironmentFile=/opt/saharlight-ux/.env
Environment=MPLCONFIGDIR=%C/mpl
ExecStart=/opt/saharlight-ux/scripts/run_bot.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
